# confs/systemd/user/swaync.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/swaync -s /home/archy/.config/MATUOUTS/matuswaync.css
